# config/systemd/mail-dovecot.service
[Unit]
Description = Mail Dovecot

[Service]
Environment=LD_LIBRARY_PATH={{ app_dir }}/dovecot/lib/dovecot
Type=simple
ExecStart={{ app_dir }}/dovecot/sbin/dovecot -F -c {{ app_dir }}/config/dovecot/dovecot.conf
NonBlocking=yes

[Install]
WantedBy = multi-user.target
